# config/systemd/user/wallpaper.service
[Unit]
Description=Auto Change Wallpapers.

[Service]
Type=simple
ExecStart=/usr/bin/change-wallpaper
